# config/guiserver.conf.template
description "GUIServer"
author "Canonical"

start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [!2345]

exec /usr/local/bin/runserver.sh
