# config/guiserver.service.template
[Unit]
Description=GUIServer

[Service]
Type=simple

ExecStart=/usr/local/bin/runserver.sh
