# config/juju-gui.conf.template
description "Juju GUI"
author "Juju GUI Peeps"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid ubuntu
exec node {juju_gui_dir}/server.js
