# config/juju-gui.conf.template
description "Juju GUI"
author "Canonical"

start on runlevel [2345]
stop on runlevel [!2345]

exec service nginx start
